# crontab/crontab
PATH=/Users/vincent/scripts:/usr/local/bin:/bin:/usr/bin
LOGFILE=/usr/local/var/log/cron.log

0 1 * * * log.sh drum-instagram $LOGFILE
0 1 * * * log.sh drum-youtube $LOGFILE
